# elasticbeanstalk-docker-dropwizard-webapp-archetype/src/main/resources/archetype-resources/src/main/config/supervisord.conf
[supervisord]
nodaemon=true

[program:service]
command=/app/bin/app.sh server /app/config-aws.yml
autostart=true
autorestart=true
